# config/systemd/user/copyq.service
[Unit]
Description=copyq

[Service]
ExecStartPre=sleep 5
ExecStart=/usr/bin/copyq

[Install]
WantedBy=Wayland.target
