# systemd/auto-snapshots.service
[Unit]
Description=Btrfs auto snapshots Service

[Service]
ExecStart = /home/chauncey/dotfiles/utils/.scripts/btrfs-autosnapshot.py -n 30 / home
