# systemd/auto-snapshots.timer

[Unit]
Description=Runs mytimer 3/D

[Timer]
OnUnitActiveSec=8h
Unit=auto-snapshots.service

[Install]
WantedBy=multi-user.target
